# system/etc/systemd/system/hypr-login-manager-autorevert.timer
[Unit]
Description=Safety timer for greetd login-manager trial

[Timer]
OnBootSec=7min
AccuracySec=30s
Unit=hypr-login-manager-autorevert.service

[Install]
WantedBy=timers.target
