# system/etc/systemd/system/hypr-login-manager-autorevert.service
[Unit]
Description=Revert greetd trial to Plasma Login Manager if unconfirmed
After=greetd.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/hypr-login-manager-autorevert
