# home/.config/systemd/user/hyprland-session.target
[Unit]
Description=Hyprland graphical session
Requires=graphical-session.target
BindsTo=graphical-session.target
Before=graphical-session.target
